# ospf/volumes/supervisord.conf
[supervisord]
nodaemon=true

#[program:true]
#command=sleep 90000
[program:zebra]
command=/usr/lib/quagga/zebra -f /etc/quagga/zebra.conf
[program:ospfd]
command=/usr/lib/quagga/ospfd -f /etc/quagga/ospfd.conf
[program:ospf6d]
command=/usr/lib/quagga/ospf6d -f /etc/quagga/ospf6d.conf
